# linux/temperature-sensor-daemon.service
[Unit]
Description=Webserver to configure temperature sensors

[Service]
ExecStart=/home/pi/temperature-app/run.sh