# linux/temperature-sensor-frontend.service
[Unit]
Description=Webserver to configure temperature sensors

[Service]
ExecStart=node /home/pi/temperature-sensor-app/server/server.js
